# root/etc/s6-overlay/s6-rc.d/init-config/run
#!/usr/bin/with-contenv bash

# Initialize basic configuration
echo "**** Initializing basic configuration... ****"

# Create necessary directories
mkdir -p /config
mkdir -p /etc/wireguard

# Set proper permissions
chmod 755 /config
chmod 755 /etc/wireguard

echo "**** Basic configuration initialized ****"
exit 0
